# elements/manila-ssh/install.d/fetch-public-ssh-keys.conf
description "Fetch public ssh key from Nova metadata service"

start on (started networking)

task

exec /usr/local/sbin/fetch-public-ssh-keys.sh
